# supervisord_sphinx.conf
[supervisord]
nodaemon=true

[program:searchd]
command=/usr/bin/searchd --config /etc/sphinx/sphinx.conf --nodetach
user=root
